# mac/pkg/postinstall
#!/bin/bash
# colourMatik .pkg postinstall — runs the full setup as the logged-in user.
# The .pkg bundles the code (staged to /Library/Application Support/colourMatik-stage),
# so no git / Xcode Command Line Tools are needed. Only Homebrew (Python + ffmpeg)
# is fetched if missing. Runs headless; the macOS Installer shows the progress bar.
set -u
LOG="/tmp/colourMatik-pkg-install.log"
exec > >(tee -a "$LOG") 2>&1
echo "== colourMatik install $(date) =="

STAGE="/Library/Application Support/colourMatik-stage"
CONSOLE_USER="$(/usr/bin/stat -f%Su /dev/console)"
[ -z "$CONSOLE_USER" ] || [ "$CONSOLE_USER" = "root" ] && CONSOLE_USER="$(/usr/bin/ls -l /dev/console | awk '{print $3}')"
USER_UID="$(/usr/bin/id -u "$CONSOLE_USER")"
USER_HOME="$(/usr/bin/dscl . -read "/Users/$CONSOLE_USER" NFSHomeDirectory | awk '{print $2}')"
DEST="$USER_HOME/colourMatik"
echo "user=$CONSOLE_USER uid=$USER_UID home=$USER_HOME"

# run a command as the logged-in user, in their GUI session
asuser() { /bin/launchctl asuser "$USER_UID" /usr/bin/sudo -u "$CONSOLE_USER" "$@"; }

# 1) place the bundled code in ~/colourMatik (user-owned)
/bin/rm -rf "$DEST"
/bin/mkdir -p "$DEST"
/usr/bin/ditto "$STAGE" "$DEST"
/usr/sbin/chown -R "$CONSOLE_USER" "$DEST"
/usr/bin/xattr -dr com.apple.quarantine "$DEST" 2>/dev/null || true

# 2) Homebrew (only if missing). Pre-create the prefix owned by the user so the
#    unattended installer needs no password.
BREW=""
[ -x /opt/homebrew/bin/brew ] && BREW=/opt/homebrew/bin/brew
[ -x /usr/local/bin/brew ]   && BREW=/usr/local/bin/brew
if [ -z "$BREW" ]; then
  echo "Installing Homebrew (unattended)..."
  /bin/mkdir -p /opt/homebrew
  /usr/sbin/chown -R "$CONSOLE_USER":admin /opt/homebrew
  asuser /bin/bash -c 'NONINTERACTIVE=1 /bin/bash -c "$(/usr/bin/curl -fsSL https://raw.githubusercontent.com/Homebrew/install/HEAD/install.sh)"' || echo "WARN: Homebrew install returned non-zero"
  [ -x /opt/homebrew/bin/brew ] && BREW=/opt/homebrew/bin/brew
fi

# 3) Python 3.11 + ffmpeg via Homebrew, then the full engine setup + AI (as user)
if [ -n "$BREW" ]; then
  echo "Installing Python 3.11 + ffmpeg..."
  asuser /bin/bash -c "\"$BREW\" install python@3.11 ffmpeg >/dev/null 2>&1 || \"$BREW\" install python@3.11 ffmpeg"
  BREW_PREFIX="$("$BREW" --prefix)"
  echo "Running engine setup (venv, deps, AI models)..."
  asuser /bin/bash -c "cd '$DEST' && PATH=\"$BREW_PREFIX/bin:\$PATH\" ./setup.sh"
  echo "Installing the Premiere panel..."
  asuser /bin/bash -c "cd '$DEST' && PATH=\"$BREW_PREFIX/bin:\$PATH\" ./install-panel.sh"
else
  echo "WARN: Homebrew unavailable; engine/AI not set up. Panel+effect still installed."
fi

# 4) native effect -> shared MediaCore (we are root here, no sudo needed)
MC="/Library/Application Support/Adobe/Common/Plug-ins/7.0/MediaCore"
if [ -d "$DEST/colourmatik-fx/colourMatik.plugin" ]; then
  /bin/mkdir -p "$MC"
  /bin/rm -rf "$MC/colourMatik.plugin"
  /usr/bin/ditto "$DEST/colourmatik-fx/colourMatik.plugin" "$MC/colourMatik.plugin"
  /usr/bin/xattr -dr com.apple.quarantine "$MC/colourMatik.plugin" 2>/dev/null || true
  echo "Effect installed."
fi

# 5) start the engine at login (LaunchAgent) + now
if [ -n "${BREW:-}" ] && [ -d "$DEST/.venv" ]; then
  BREW_PREFIX="$("$BREW" --prefix)"
  PLIST="$USER_HOME/Library/LaunchAgents/com.colourmatik.engine.plist"
  asuser /bin/mkdir -p "$USER_HOME/Library/LaunchAgents" "$USER_HOME/Library/Logs"
  /bin/cat > "$PLIST" <<PL
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0"><dict>
  <key>Label</key><string>com.colourmatik.engine</string>
  <key>ProgramArguments</key><array>
    <string>$DEST/.venv/bin/python</string><string>-u</string><string>-m</string><string>colourmatik.webapp</string>
  </array>
  <key>WorkingDirectory</key><string>$DEST</string>
  <key>EnvironmentVariables</key><dict><key>PATH</key><string>$BREW_PREFIX/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string></dict>
  <key>LimitLoadToSessionType</key><string>Aqua</string>
  <key>RunAtLoad</key><true/><key>KeepAlive</key><true/>
  <key>ProcessType</key><string>Interactive</string><key>ThrottleInterval</key><integer>5</integer>
  <key>StandardOutPath</key><string>$USER_HOME/Library/Logs/colourmatik-engine.log</string>
  <key>StandardErrorPath</key><string>$USER_HOME/Library/Logs/colourmatik-engine.log</string>
</dict></plist>
PL
  /usr/sbin/chown "$CONSOLE_USER" "$PLIST"
  asuser /bin/launchctl bootout "gui/$USER_UID/com.colourmatik.engine" 2>/dev/null || true
  asuser /bin/launchctl bootstrap "gui/$USER_UID" "$PLIST" 2>/dev/null || \
  asuser /bin/launchctl load "$PLIST" 2>/dev/null || true
  echo "Engine autostart configured."
fi

echo "== colourMatik install done $(date) =="
exit 0
